# srv/salt/ceph/iperf/systemd-iperf.service
[Unit]
Description=iperf server

[Service]
ExecStart=/usr/bin/iperf3 -s

[Install]
WantedBy=default.target
